# deploy/twitch-drops-bot.service
[Unit]
Description=Twitch Drops Bot — Telegram Commands
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projects/twitch-drops-bot
EnvironmentFile=%h/projects/twitch-drops-bot/.env
ExecStart=%h/projects/twitch-drops-bot/venv/bin/python bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
