# deploy/twitch-drops-check.timer
[Unit]
Description=Twitch Drops Bot — Daily Drop Check Timer

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
